# install_mercury/mercury.service.in
[Unit]
Description=mercury packet metadata capture and analysis
Documentation=man:mercury(1)
After=network.target

[Service]
Type=simple
WorkingDirectory=@prefix@/var/mercury/
ExecStart=@prefix@/bin/mercury --config /etc/mercury/mercury.cfg
StandardOutput=file:/var/log/mercury.log
StandardError=file:/var/log/mercury.log

[Install]
WantedBy=multi-user.target
